# vps/gb_demo_gspc.conf
[program:gb-demo-gspc]
command=/home/gb-brain/GB-BRAIN/.venv/bin/python -m runtime.demo_runner --broker alpaca --symbol ^GSPC --tf 1h --strategy gspc_gem_1
directory=/home/gb-brain/GB-BRAIN
user=root
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/var/log/gb-brain/demo_gspc.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/gb-brain/demo_gspc_err.log
stderr_logfile_maxbytes=5MB
environment=PATH="/home/gb-brain/GB-BRAIN/.venv/bin:%(ENV_PATH)s"
