# vps/supervisor_observer.conf
; GB-BRAIN — Supervisor config for live_observer
; ================================================
; Manages the GB-BRAIN live market observer process.
; Monitors active lanes, scores signals, and logs
; observed signals to SQLite.
;
; Install: sudo cp this /etc/supervisor/conf.d/gb-observer.conf
;          sudo supervisorctl reread && sudo supervisorctl update

[program:gb-observer]

; --- Command & working directory ---
command=/opt/gb-brain/venv/bin/python /opt/gb-brain/runtime/live_observer.py --mode paper --lane GB-INDICES
directory=/opt/gb-brain

; --- Process lifecycle ---
autostart=true
autorestart=true
startsecs=5
startretries=5
stopwaitsecs=30
killasgroup=true
stopasgroup=true

; --- User ---
user=gb

; --- Environment ---
; Injects the venv bin directory so any subprocess spawned by the observer
; also resolves the correct Python / pip.
environment=PATH="/opt/gb-brain/venv/bin:%(ENV_PATH)s",PYTHONUNBUFFERED="1"

; --- Logging ---
stdout_logfile=/var/log/gb-brain/observer.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

stderr_logfile=/var/log/gb-brain/observer_err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3

; Merge stdout+stderr into one stream for easier tailing (optional)
; redirect_stderr=true

; --- Priority (lower = started earlier) ---
priority=100
